# supervisor/elasticSearch.conf
[program:elasticSearch]
command=/usr/bin/start-elasticSearch.sh
autostart=true
priority=15